# systemd/nordvpn-gateway.service
systemd-fil 


Åpne filen for redigering:

sudo nano /etc/systemd/system/nordvpn-gateway.service

Slett ALT innholdet som er der nå.
Lim inn den enkle, rene versjonen:

[Unit]
Description=NordVPN Gateway Service
# Venter på at nettverket skal være "online" og at pihole skal være klar
After=network-online.target pihole-FTL.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/nordvpn-gateway.env
User=root

# VIKTIG: Denne linjen er ny! Den venter til den kan pinge gatewayen.
ExecStartPre=/bin/bash -c 'echo "Venter på at UCG-Fiber (192.168.1.1) skal svare..."; while ! ping -c 1 -W 2 192.168.1.1 &>/dev/null; do sleep 5; done; echo "UCG-Fiber svarer, starter hovedscript."'

# Hovedscriptet som starter ETTER at ExecStartPre er vellykket
ExecStart=/usr/local/bin/nordvpn-gateway.sh

Restart=always
RestartSec=30 # Øker ventetiden litt mellom restart-forsøk

StandardOutput=file:/var/log/nordvpn-gateway.log
StandardError=file:/var/log/nordvpn-gateway.log

[Install]
WantedBy=multi-user.target